# kanata/.config/systemd/user/kanata.service
[Unit]
Description=Improve keyboard comfort and usability with advanced customization

[Service]
Type=simple
ExecStart=kanata -c "$HOME/.config/kanata/config.kbd"

[Install]
WantedBy=default.target
